# ansible/roles/void-core/templates/rc.local.j2
#!/bin/sh
# rc.local - Local startup script
# Generated by Ansible

{% if void_rc_local is defined and void_rc_local %}
{{ void_rc_local }}
{% else %}
# Add your custom startup commands here
{% endif %}

# Example commands:
# echo "Custom startup script executed" > /var/log/rc.local.log
# mount -t tmpfs tmpfs /tmp -o size=1G
# /usr/local/bin/my-startup-script.sh

exit 0
